# programming-language/go/Building-RESTful-Web-Services-with-Go/Chapter01/goproject.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:romanserver]
# hardcode in windows 10 WSL2 ubuntu22.04
command=/root/go/bin/romanserver
autostart=true
autorestart=true
redirect_stderr=true
